# checker.conf
#!upstart

description "gamertagchecekr.com server"
author "Jackson Owens"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on shutdown

respawn
respawn limit 5 60

script
  exec sudo uwsgi --ini /var/www/checker/checker.ini
end script
